# etc/profile.d/xsdb-variables.sh
export XILINX_VITIS=/usr/local/xilinx_vitis
export PATH=${PATH}:${XILINX_VITIS}
export TCLLIBPATH=/usr/local/xilinx_vitis
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/local/lib
export LC_ALL='en_US.UTF-8'
